# scripts/supervisord.conf
[supervisord]
nodaemon=false
logfile=/workspace/shared/minicluster/logs/supervisord.log
pidfile=/workspace/shared/minicluster/supervisord.pid

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[program:redis]
command=taskset -c 0-55 redis-server --port 6379 --save '' --appendonly no --protected-mode yes --bind 127.0.0.1
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/redis.log
stderr_logfile=/workspace/shared/minicluster/logs/redis.err

[program:payments]
command=taskset -c 0-55 python3 app.py
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/payments.log
stderr_logfile=/workspace/shared/minicluster/logs/payments.err
environment=SERVICE_NAME="payments",PORT="7001",NODE="node-1"

[program:auth]
command=taskset -c 0-55 python3 app.py
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/auth.log
stderr_logfile=/workspace/shared/minicluster/logs/auth.err
environment=SERVICE_NAME="auth",PORT="7002",NODE="node-1"

[program:checkout]
command=taskset -c 56-111 python3 app.py
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/checkout.log
stderr_logfile=/workspace/shared/minicluster/logs/checkout.err
environment=SERVICE_NAME="checkout",PORT="7003",NODE="node-2"

[program:fraud]
command=taskset -c 56-111 python3 app.py
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/fraud.log
stderr_logfile=/workspace/shared/minicluster/logs/fraud.err
environment=SERVICE_NAME="fraud",PORT="7004",NODE="node-2"

[program:loadgen]
command=taskset -c 108-111 python3 loadgen.py
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/loadgen.log
stderr_logfile=/workspace/shared/minicluster/logs/loadgen.err

[program:collector]
command=taskset -c 108-111 python3 collector.py
directory=/workspace/shared/minicluster
autorestart=true
startsecs=2
stopwaitsecs=5
stdout_logfile=/workspace/shared/minicluster/logs/collector.log
stderr_logfile=/workspace/shared/minicluster/logs/collector.err
